# cluster-playbook/roles/influxdb-load-test/files/influxdb-load-test.service
[Unit]
Description=Load test Influxdb

[Service]
Type=oneshot
ExecStart=/usr/local/bin/load-test-influxdb
