# cluster-playbook/roles/influxdb-load-test/files/influxdb-load-test.timer
[Unit]
Description=Load test InfluxDB every 5 minutes

[Timer]
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
